# asterisk/rootfs/etc/services.d/stdin/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Asterisk
# Starts the STDIN service
# ==============================================================================

bashio::log.info 'Starting the Home Assistant STDIN service...'

while read -r input; do
    input=$(bashio::jq "${input}" '.')
    bashio::log.info "Read STDIN command: $input"

    case "$input" in

        "reload_pjsip")
            cp /config/asterisk/pjsip_custom.conf /etc/asterisk/pjsip_custom.conf
            cp /config/asterisk/pjsip.conf /etc/asterisk/pjsip.conf
            asterisk -x module reload res_pjsip
        ;;

        "reload_dialplan")
            cp /config/asterisk/extensions.conf /etc/asterisk/extensions.conf
            asterisk -x dialplan reload
        ;;

        *)
            bashio::log.warning "Unknown STDIN command: $input"
        ;;

    esac
done